# init/jemaos-power-daemon.conf
# Copyright (c) 2025 Jema Technology. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Start the JemaOS power daemon"

start on started powerd
stop on stopping boot-services

oom score -100

script
  # Execute the power daemon and redirect logs to a temporary file
  exec /usr/sbin/power_daemon > /tmp/jemaos_powerd.log
end script